# ztsync.service
[Unit]
Description=1ClickSync Flask Application
After=network.target postgresql.service redis.service

[Service]
User=zmarin
WorkingDirectory=/home/zmarin/Projects/ZTsync_python
Environment="PATH=/home/zmarin/Projects/ZTsync_python/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/zmarin/Projects/ZTsync_python"
Environment="FLASK_APP=/home/zmarin/Projects/ZTsync_python/app/__init__.py"
Environment="FLASK_ENV=development"
Environment="FLASK_DEBUG=1"
EnvironmentFile=/home/zmarin/Projects/ZTsync_python/.env

ExecStart=/bin/bash -c 'source /home/zmarin/Projects/ZTsync_python/venv/bin/activate && /home/zmarin/Projects/ZTsync_python/start-server.sh'

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
